# com.hermes.chrome-cdp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.chrome-cdp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{REPO}}/scripts/hermes-chrome-cdp.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>HERMES_CDP_PORT</key>
        <string>9222</string>
        <key>HERMES_CDP_BIND</key>
        <string>127.0.0.1</string>
        <key>HERMES_CDP_RECLAIM_SQUAT</key>
        <string>1</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/hermes-chrome-cdp.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/hermes-chrome-cdp.log</string>
</dict>
</plist>
