# com.igor.thumbgate-aeo-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.igor.thumbgate-aeo-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{NODE}}</string>
    <string>{{RUNTIME}}/tools/thumbgate-aeo-monitor.js</string>
    <string>--execute</string>
    <string>--write</string>
    <string>--json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>{{RUNTIME}}</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>PARALLEL_CLI_BIN</key>
    <string>{{PARALLEL_CLI}}</string>
    <key>PATH</key>
    <string>{{PATH}}</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>LowPriorityIO</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>{{HOME}}/Library/Logs/mac-yolo/thumbgate-aeo-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>{{HOME}}/Library/Logs/mac-yolo/thumbgate-aeo-monitor.error.log</string>
</dict>
</plist>
