# com.igor.memory-pressure-guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.igor.memory-pressure-guardian</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.local/bin/memory-pressure-guardian.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>MEMORY_GUARD_SHED_HERMES_GATEWAY</key>
        <string>1</string>
        <key>MEMORY_GUARD_HERMES_GATEWAY_LABEL</key>
        <string>ai.hermes.gateway</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/memory-pressure-guardian.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/memory-pressure-guardian.launchd.log</string>
</dict>
</plist>
