# com.igor.hermes-mobile-pair-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.igor.hermes-mobile-pair-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{NODE}}</string>
        <string>{{REPO}}/tools/hermes-mobile-pair.js</string>
        <string>--server-only</string>
    </array>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>ThrottleInterval</key><integer>10</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>{{HOME}}</string>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key><string>{{HOME}}/Library/Logs/hermes-mobile-pair-server.log</string>
    <key>StandardErrorPath</key><string>{{HOME}}/Library/Logs/hermes-mobile-pair-server.log</string>
</dict>
</plist>
