# com.igor.hermes-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.igor.hermes-dashboard</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/igorganapolsky/.local/bin/node</string>
		<string>/Users/igorganapolsky/workspace/git/igor/mac-yolo-safeguards/tools/hermes-dashboard.js</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/hermes-dashboard.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/hermes-dashboard.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HERMES_DASH_PORT</key>
		<string>8787</string>
		<key>HERMES_DASH_HOST</key>
		<string>127.0.0.1</string>
	</dict>
</dict>
</plist>
